# quotes_service/package/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable quotes
systemctl start quotes
